# .config/systemd/user/fehbg.service
[Unit]
Description=Restore background
Requires=wm.target

[Service]
ExecStart=%h/.fehbg

[Install]
WantedBy=wm.target
