# .config/systemd/user/fehbg.service
[Unit]
Description=Restore background
Requires=x11.target

[Service]
ExecStart=%h/.fehbg

[Install]
WantedBy=x11.target
