# .config/systemd/user/kbd-hotplug-udev.service
[Unit]
Description=udev xmodmap monitor
After=x11.target

[Service]
ExecStart=%h/.local/lib/X11/kbd-hotplug-udev

[Install]
WantedBy=x11.target
