# .config/systemd/user/spotify-listener.service
[Unit]
Description=Watch for spotify open for polybar controls

[Service]
ExecStart=%h/.local/lib/polybar/spotify_listener.py

