# .config/systemd/user/wm.target
[Unit]
Description=Window manager has started
