# .config/systemd/user/xbanish.service
[Unit]
Description=Hide cursor when key is pressed
Requires=x11.target

[Service]
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=x11.target
