# .config/systemd/user/xidlehook.service
[Unit]
Description=Automatic screen locker
Requires=pulseaudio.service
Requires=wm.target

[Service]
ExecStart=%h/.local/bin/screensaver.sh
Restart=always

[Install]
WantedBy=wm.target
