# .config/systemd/user/xidlehook.service
[Unit]
Description=Automatic screen locker
After=pulseaudio
Requires=wm.target

[Service]
ExecStart=/usr/bin/xidlehook --not-when-audio --timer 300 %h/.local/bin/lock.sh ''
Restart=always

[Install]
WantedBy=wm.target
