# .config/systemd/user/xsetroot.service
[Unit]
Description="Display info on dwm bar"
Requires=dwm.target

[Service]
ExecStart=%h/.local/lib/dwm/xsetroot.sh

[Install]
WantedBy=dwm.target
